# com.auotam.email-automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.auotam.email-automation</string>

  <key>WorkingDirectory</key>
  <string>/Users/govindchauhan/Documents/auotam-lead-genertion</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>/Users/govindchauhan/Documents/auotam-lead-genertion/run_prod.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/govindchauhan/Documents/auotam-lead-genertion/output/email/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/govindchauhan/Documents/auotam-lead-genertion/output/email/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
